# src/openwebui/openwebui.conf
[program:openwebui]
redirect_stderr=true
autorestart=true
command=/opt/start-openwebui.sh