# src/rstats/rstudio.conf
[program:rstudio]
autostart=true
command=/opt/start-rstudio.sh