# src/jupyterhub/install-jupyterhub.conf
# supervisord that runs a script that installs JupyterHub
# into the Kubernetes cluster.

[program:jupyterhub]
redirect_stderr=true
command=/jupyterhub/install-jupyterhub.sh
autorestart=false