# src/hpc/slurm-supervisor.conf
[program:slurm]
command=sudo /etc/slurm/start-slurm.sh
autostart=true
