# src/cocalc/supervisor/conf.d/code-server.conf
[program:code-server]
directory=/home/user
user=user
environment=HOME="/home/user",USER="user",PORT="8123"
command=/cocalc/apps/code-server.sh
autostart=false
autorestart=true
