# src/cocalc/supervisor/conf.d/cron.conf
[program:cron]
user=root
command=/usr/sbin/cron -f
autostart=true
