# src/cocalc/supervisor/conf.d/cron.conf
[program:cron]
command=sudo /usr/sbin/cron -f
autostart=true
