# src/cocalc/supervisor/conf.d/cron.conf
[program:cron]
command=/usr/sbin/cron -f -s
autostart=true
