# src/cocalc/supervisor/conf.d/jupyterlab.conf
[program:jupyterlab]
directory=/home/user
user=user
environment=HOME="/home/user",PORT="8888"
command=/cocalc/apps/jupyterlab.sh
autostart=false
autorestart=true
