# src/cocalc/supervisor/conf.d/pluto.conf
[program:pluto]
directory=/home/user
user=user
environment=HOME="/home/user",USER="user",PORT="8124"
command=/cocalc/apps/pluto.sh
autostart=false
autorestart=true
