# src/cocalc/supervisor/conf.d/xpra.conf
[program:xpra]
directory=/home/user
user=user
environment=HOME="/home/user",USER="user",PORT="10000"
command=/cocalc/apps/xpra.sh
autostart=false
autorestart=true
